# target/linux/qualcommax/ipq807x/base-files/etc/hotplug.d/usb/90-nbg7815-led
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0

seconds="$(uci -q get nbg7815_led.main.usb_notify_seconds)"
[ -n "$seconds" ] || seconds=2

/sbin/nbg7815-led-state set usb_notify
(
	sleep "$seconds"
	/sbin/nbg7815-led-state clear usb_notify
) &
