# target/linux/qualcommax/ipq807x/base-files/etc/hotplug.d/iface/90-nbg7815-led
#!/bin/sh

[ "$INTERFACE" = "wan" ] || exit 0

/sbin/nbg7815-led-state apply
