# target/linux/qualcommax/ipq807x/base-files/etc/uci-defaults/14_nbg7815_led
#!/bin/sh

[ "$(board_name)" = "zyxel,nbg7815" ] || exit 0

/etc/init.d/nbg7815-led enable
/etc/init.d/nbg7815-led restart >/dev/null 2>&1

exit 0
